# recipes-device-drivers/gpioservice/files/gpio_buzzer_driver_boot.service
[Unit]
Description=GPIO Driver Initialization at Bootup

[Service]
ExecStart=/bin/sh -c '/lib/modules/$(uname -r)/aesdgpio_load'

[Install]
WantedBy=multi-user.target
